# voc.service
[Unit]
Description=VOC vocabulary game – Node.js API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/voc
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
Environment=PORT=3000
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
